# kiosk.conf
#!/bin/bash
#This file is placed in /etc/init/kiosk.conf
start on (filesystem and stopped udevtrigger)
stop on runlevel [06]

console output
emits starting-x

respawn

exec sudo -u kiosk startx /etc/X11/Xsession /opt/kiosk.sh --